# aibox-hermes/provision/aibox-webui@.service
# Service systemd — une interface web hermes-webui PRIVÉE par utilisateur (derrière Caddy).
# Instance %i = nom d'utilisateur (ex: aibox-webui@andre). Remplace aibox-dash@.
#
# Bind 127.0.0.1 (loopback) → seul Caddy l'atteint ; jamais exposé directement.
#
# GABARIT : les jetons __AIBOX_OWNER__, __OWNER_HOME__ et __AIBOX_REPO__ sont
# substitués À L'INSTALLATION par setup-portal.sh (sed), à partir de l'owner réel
# et du chemin réel du repo. Sur xefia c'était clikinfo + /home/clikinfo/BoxIA ;
# sur toute autre machine (owner ≠ clikinfo, repo ≠ ~/BoxIA) les chemins suivent,
# sinon l'extension white-label servait un 404.
#
# EnvironmentFile par user — __OWNER_HOME__/aibox/webui/<user>.env :
#     HERMES_HOME=__OWNER_HOME__/aibox/companies/<co>/users/<user>/hermes
#     HERMES_WEBUI_PORT=9130        # un port distinct par user (9130, 9131, …)
#     HERMES_WEBUI_EXTENSION_SCRIPT_URLS=/extensions/aibox.js?role=admin|client
[Unit]
Description=AI Box - Interface web utilisateur (%i)
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=__AIBOX_OWNER__
EnvironmentFile=__OWNER_HOME__/aibox/webui/%i.env
Environment=PATH=__OWNER_HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HERMES_WEBUI_AGENT_DIR=__OWNER_HOME__/hermes-agent
Environment=HERMES_WEBUI_PYTHON=__OWNER_HOME__/hermes-agent/venv/bin/python
Environment=HERMES_WEBUI_HOST=127.0.0.1
Environment=HERMES_WEBUI_EXTENSION_DIR=__AIBOX_REPO__/aibox-hermes/webui-ext/static
Environment=HERMES_WEBUI_EXTENSION_STYLESHEET_URLS=/extensions/aibox.css
# HERMES_WEBUI_EXTENSION_SCRIPT_URLS est défini PAR USER dans l'EnvironmentFile
# (/extensions/aibox.js?role=admin|client) → pilote les droits selon le rôle.
WorkingDirectory=__OWNER_HOME__/hermes-webui
ExecStart=/usr/bin/python3 __OWNER_HOME__/hermes-webui/bootstrap.py --foreground --no-browser
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
