# services/setup/aibox-mdns-aliases.service
[Unit]
Description=Publish AI Box mDNS aliases (aibox.local + sous-domaines)
Documentation=https://github.com/Tobaz34/BoxIA/blob/main/docs/ACCESS.md
After=avahi-daemon.service network-online.target
Requires=avahi-daemon.service
Wants=network-online.target

[Service]
Type=simple
# Lit DOMAIN_PREFIX écrit par le wizard dans /srv/ai-stack/.env. Le `-`
# rend le fichier optionnel (le service démarre quand même au 1er boot,
# avant que le wizard ait écrit .env). Quand .env change, redémarrer
# manuellement : sudo systemctl restart aibox-mdns-aliases.service
EnvironmentFile=-/srv/ai-stack/.env
ExecStart=/usr/local/bin/aibox-mdns-publish.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
